# lib/newshound/scripts/cron.d/newshound_fetch_cron
*/1 * * * * root /bin/bash /opt/newshound/bin/newshound-fetch

