# newshound/scripts/cron.d/newshound_fetch_cron
*/1 * * * * root /bin/bash /opt/nyt/bin/newshound-fetch

